# templates/vpn-web.service
# templates/vpn-web.service
[Unit]
Description=OpenVPN Web Interface
After=network.target
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ovpn-client-web
EnvironmentFile=/opt/ovpn-client-web/.env
ExecStart=/usr/bin/make run
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target